# alertmanager.service
[Unit]
Description=Alertmanager Service
After=network.target
[Service]
User=prometheus
Group=prometheus
Type=simple
ExecStart=/usr/local/bin/alertmanager --config.file=/etc/prometheus/alertmanager.yml --storage.path=/etc/prometheus/alertmanager/data
Restart=always
[Install]
WantedBy=multi-user.target
